# itamae/cookbooks/cms/files/etc/systemd/system/cms-logservice.service
[Unit]
Description=CMS LogService

[Service]
User=cmsuser
Group=cmsuser
RuntimeDirectory=cms

ExecStart=/usr/bin/cmsLogService

[Install]
WantedBy=multi-user.target
