# itamae/cookbooks/cms/files/etc/systemd/system/cms-checker.service
[Unit]
Description=CMS Checker

[Service]
User=cmsuser
Group=cmsuser
RuntimeDirectory=cms

ExecStart=/usr/bin/cmsChecker

[Install]
WantedBy=multi-user.target
