# itamae/cookbooks/cms/files/etc/systemd/system/cms-proxyservice.service
[Unit]
Description=CMS ProxyService

[Service]
User=cmsuser
Group=cmsuser
RuntimeDirectory=cms

EnvironmentFile=-/etc/ioi-contest-id.env
ExecStart=/usr/bin/cmsProxyService -c $CMS_CONTEST_ID

[Install]
WantedBy=multi-user.target
