# itamae/cookbooks/cms/files/etc/systemd/system/cms.target
[Unit]
Description=CMS system

[Install]
WantedBy=multi-user.target
