# itamae/cookbooks/cms/files/etc/systemd/system/ioi-update-cms-config.timer
[Unit]
Description=periodically update CMS config

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min
RandomizedDelaySec=15

[Install]
WantedBy=timers.target
