# server/cookbooks/cms/files/etc/systemd/system/cms-resourceservice.service
[Unit]
Description=CMS ResourceService

[Service]
User=cmsuser
Group=cmsuser
RuntimeDirectory=cms

ExecStart=/usr/bin/cmsResourceService

[Install]
WantedBy=multi-user.target
