# server/cookbooks/cms/files/etc/systemd/system/cms-scoringservice.service
[Unit]
Description=CMS ScoringService

[Service]
User=cmsuser
Group=cmsuser
RuntimeDirectory=cms

ExecStart=/usr/bin/cmsScoringService

[Install]
WantedBy=multi-user.target
