# itamae/roles/dhcp/files/etc/systemd/system/kea-dhcp4-server.service.d/directory.conf
[Service]
RuntimeDirectory=kea
StateDirectory=kea
