# itamae/cookbooks/codedeploy-agent/files/etc/systemd/system/ioi-ensure-onpremise-codedeploy-credentials.timer
[Unit]
Description=Ensure AWS credentials for CodeDeploy onpremises

[Timer]
OnBootSec=1min
OnUnitActiveSec=4hr

[Install]
WantedBy=timers.target
