# itamae/roles/dhcp/files/etc/systemd/system/ioi-upload-leases.service
[Unit]
Description=Upload Kea leases
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/ioi-upload-leases
