# itamae/cookbooks/ioi-set-hostname/files/etc/systemd/system/ioi-set-hostname.service
[Unit]
Description=Set hostname based on MAC address

[Service]
Type=oneshot
ExecStart=/usr/sbin/ioi-set-hostname

[Install]
WantedBy=multi-user.target
