# itamae/cookbooks/prometheus-cups/files/etc/systemd/system/prometheus-cups.timer
[Unit]
Description=Collect CUPS metrics for prometheus
After=cups.service

[Timer]
OnUnitInactiveSec=60s

[Install]
WantedBy=timers.target
