# itamae/cookbooks/prometheus-usb/files/etc/systemd/system/prometheus-usb.timer
[Unit]
Description=Collect connected USB devices for prometheus

[Timer]
OnBootSec=60s
OnUnitInactiveSec=30s

[Install]
WantedBy=timers.target
